# image/etc/service/gitlab-ci-multi-runner/run
#!/bin/bash

gitlab-ci-multi-runner run
